# raspberry-pi/lumiere.upstart
description "Manage Lumiere"
author "Alan Palazzolo"

start on (started networking)
stop on shutdown

# Restart the process if it dies with a signal
# or exit code not given by the 'normal exit' stanza.
respawn

# Give up if restart occurs 10 times in 90 seconds.
respawn limit 10 90

script
  export HOME="/home/pi"
  exec python /home/pi/lumiere/raspberry-pi/lumiere.py
end script